# data/shootbacks.ini
[program:shootbacks]
command=python3 /opt/shootback-master/master.py -m 0.0.0.0:10022 -c 0.0.0.0:10023 -k yourpassword
autorestart=true
user=nobody
